# supervisord.conf
[supervisord]
nodaemon=true

[program:poormans-ddns]
command=/pmddns/poormans-ddns -c /pmddns/config/pmddns.toml
